# eredes-token/rootfs/etc/services.d/eredes-token/run
#!/usr/bin/with-contenv bashio

bashio::log.info "==================================="
bashio::log.info "  E-Redes Token Refresher Starting"
bashio::log.info "==================================="

# Start Xvfb
bashio::log.info "Starting Xvfb..."
Xvfb :99 -screen 0 1280x720x24 &
sleep 3

export DISPLAY=:99

# Start x11vnc
bashio::log.info "Starting x11vnc..."
x11vnc -display :99 -forever -shared -nopw -rfbport 5901 -q &
sleep 2

# Start noVNC on port 6081
bashio::log.info "Starting noVNC on port 6081..."
/opt/novnc/utils/novnc_proxy --vnc localhost:5901 --listen 6081 &
sleep 2

bashio::log.info "==================================="
bashio::log.info "  noVNC: http://<your-ip>:6081"
bashio::log.info "==================================="

# Start refresher
bashio::log.info "Starting E-Redes Token Refresher..."
exec python3 /app/eredes_token.py
